# templates/default/cube-evaluator.conf.erb
description "Cube Evaluator"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

chdir <%= node.cube.install_dir %>

exec sudo -u node /usr/bin/nodejs node_modules/cube/bin/evaluator.js > <%= node[:cube][:log_dir] %>/evaluator.log 2>&1
